# .config/systemd/user/fetchcord.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/home/jacob/.local/bin/fetchcord --nohost --noshell --nohardware

